# systemd/wilt-frontend.service
[Unit]
Description=Wilt Annotation Tool Frontend
After=network.target wilt-backend.service

[Service]
User=happyai2023
WorkingDirectory=/home/happyai2023/20251208_WiltAnnotationTool/frontend
Environment="PATH=/home/happyai2023/anaconda3/envs/WiltAnnotation/bin:/usr/bin:/bin"
ExecStart=/home/happyai2023/anaconda3/envs/WiltAnnotation/bin/npm run dev -- --host
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
